# pycec.service
[Unit]
Description=pycec

[Service]
ExecStart=/usr/local/bin/pycec

[Install]
WantedBy=multi-user.target
